# lib/systemd/system/camdelete.service
[Unit]
Description=record dashcam

[Service]
ExecStart=/bin/bash /home/pi/CamProj/scripts/shell/delete.sh

[Install]
WantedBy=multi-user.target
