# etc/rc.local
#!/bin/sh -e
#
# rc.local
# Print the IP address
_IP=$(hostname -I) || true
if [ "$_IP" ]; then
  printf "My IP address is %s\n" "$_IP"
fi
python /home/pi/.off_button.py &
exit 0
